# build/overlay/etc/local.d/chromium_startup.start
#!/bin/sh

chromium --no-sandbox --disable-infobars --kiosk /index.html
